# target/linux/layerscape/base-files/etc/uci-defaults/12_mono-gateway-flowtable-reload
#!/bin/sh

. /lib/functions/system.sh

case "$(board_name)" in
mono,gateway-dk|\
mono,gateway-dk-sdboot)
		;;
*)
	exit 0
	;;
esac

changed=0

case "$(uci -q get firewall.@defaults[0].lower_device_flowtable)" in
'')
	uci set firewall.@defaults[0].lower_device_flowtable='1'
	changed=1
	;;
esac

case "$(uci -q get firewall.@defaults[0].lower_device_flowtable_reload)" in
'')
	uci set firewall.@defaults[0].lower_device_flowtable_reload='1'
	changed=1
	;;
esac

[ "$changed" -eq 0 ] || uci commit firewall

exit 0
